# com.diogo.devlab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.diogo.devlab</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/diogoangelim/devlab/start-devlab.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>3000</string>
    <key>SERVICE_NAME</key>
    <string>Diogo Dev Lab</string>
    <key>NODE_BIN</key>
    <string>/opt/homebrew/bin/node</string>
    <key>ENABLE_HOTSPOT</key>
    <string>0</string>
    <key>BIND_ADDR</key>
    <string>127.0.0.1</string>
    <key>TOKEN_FILE</key>
    <string>/Users/diogoangelim/devlab/config/token</string>
    <key>SHARE_SERVICE</key>
    <string>0</string>
    <key>ENABLE_TLS</key>
    <string>1</string>
    <key>TLS_KEY_FILE</key>
    <string>/Users/diogoangelim/devlab/config/server.key.pem</string>
    <key>TLS_CERT_FILE</key>
    <string>/Users/diogoangelim/devlab/config/server.cert.pem</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/diogoangelim/devlab/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/diogoangelim/devlab/logs/launchd.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/diogoangelim/devlab</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
